# backend/services/yata-backend.service
[Unit]
Description=YATA Backend Service
After=network.target

[Service]
Type=simple
# 注意: 此配置仅适用于用户级服务 (systemctl --user)
# 如需系统级服务，请使用 yata-backend-system.service
WorkingDirectory=%h/YATA/backend
Environment="PATH=%h/YATA/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=%h/YATA/backend/.venv/bin/python src/run_service.py
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

