# backend/systemd-config/frontend/yata-frontend.service
[Unit]
Description=YATA Frontend Service (Production)
After=network.target

[Service]
Type=simple
# 注意: 此配置仅适用于用户级服务 (systemctl --user)
# 生产模式: 使用 npm start (需要先运行 npm run build)
WorkingDirectory=%h/YATA/frontend
Environment="NODE_ENV=production"
Environment="PORT=3000"
# 通过 bash 加载 nvm 环境后执行 npm (支持 nvm 管理的 Node.js)
ExecStart=/bin/bash -c 'source %h/.nvm/nvm.sh && npm start'
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
