# backend/systemd-config/frontend/yata-frontend-dev.service
[Unit]
Description=YATA Frontend Service (Development)
After=network.target

[Service]
Type=simple
# 注意: 此配置仅适用于用户级服务 (systemctl --user)
# 开发模式: 使用 npm run dev (支持热重载)
WorkingDirectory=%h/YATA/frontend
Environment="PATH=%h/YATA/frontend/node_modules/.bin:/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=development"
Environment="PORT=3000"
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
