# backend/systemd-config/RAG/yata-rag.service
[Unit]
Description=YATA RAG Service
After=network.target

[Service]
Type=simple
# 注意: 此配置仅适用于用户级服务 (systemctl --user)
WorkingDirectory=%h/YATA/algorithms/RAG_chroma
Environment="PATH=%h/YATA/algorithms/RAG_chroma/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=%h/YATA/algorithms/RAG_chroma/.venv/bin/uvicorn api_server:app --host 0.0.0.0 --port 8001
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
